# grocy_recipes/rootfs/etc/s6-overlay/s6-rc.d/recipe-backend/run
#!/command/with-contenv bashio

bashio::log.info "Starting recipe backend on port 8100..."

# ── Storage URL auto-detect ────────────────────────────────────────────────
STORAGE_URL=""

if bashio::config.has_value 'storage_url'; then
    STORAGE_URL=$(bashio::config 'storage_url')
    bashio::log.info "Storage URL (from config): ${STORAGE_URL}"
else
    bashio::log.info "Attempting auto-detect of HA-Storage addon..."

    MY_HOSTNAME="${HOSTNAME:-}"
    if [ -z "${MY_HOSTNAME}" ]; then
        MY_HOSTNAME=$(hostname 2>/dev/null) || true
    fi
    if [ -z "${MY_HOSTNAME}" ] && [ -f /etc/hostname ]; then
        MY_HOSTNAME=$(cat /etc/hostname 2>/dev/null) || true
    fi

    if [ -n "${MY_HOSTNAME}" ]; then
        MY_HOSTNAME_NORM=$(echo "${MY_HOSTNAME}" | tr '_' '-')
        REPO_PREFIX="${MY_HOSTNAME_NORM%%-grocy-recipes}"
        if [ "${REPO_PREFIX}" != "${MY_HOSTNAME_NORM}" ] && [ -n "${REPO_PREFIX}" ]; then
            for CANDIDATE in "${REPO_PREFIX}-ha-storage" "${REPO_PREFIX}_ha_storage"; do
                if curl -sf --max-time 3 "http://${CANDIDATE}:8099/api/products" >/dev/null 2>&1; then
                    STORAGE_URL="http://${CANDIDATE}:8099"
                    bashio::log.info "Storage URL (auto-detected via hostname): ${STORAGE_URL}"
                    break
                fi
            done
        fi
    fi

    # Supervisor API fallback
    if [ -z "${STORAGE_URL}" ] && [ -n "${SUPERVISOR_TOKEN}" ]; then
        STORAGE_SLUG=""
        STORAGE_SLUG=$(curl -sf --max-time 5 \
            -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
            http://supervisor/addons \
            | jq -r '
                [.data.addons[]?
                 | select(.slug | test("ha.storage"; "i"))
                ][0].slug // empty
            ' \
        ) || true

        if [ -n "${STORAGE_SLUG}" ]; then
            STORAGE_IP=""
            STORAGE_IP=$(curl -sf --max-time 5 \
                -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
                "http://supervisor/addons/${STORAGE_SLUG}/info" \
                | jq -r '.data.ip_address // empty' \
            ) || true

            if [ -n "${STORAGE_IP}" ]; then
                STORAGE_URL="http://${STORAGE_IP}:8099"
                bashio::log.info "Storage URL (auto-detected via API): ${STORAGE_URL}"
            elif curl -sf --max-time 3 "http://${STORAGE_SLUG}:8099/api/products" >/dev/null 2>&1; then
                STORAGE_URL="http://${STORAGE_SLUG}:8099"
                bashio::log.info "Storage URL (auto-detected via slug): ${STORAGE_URL}"
            fi
        fi
    fi

    if [ -z "${STORAGE_URL}" ]; then
        bashio::log.warning "HA-Storage addon not found — will retry every 5 s (up to 20 attempts)."
        for _attempt in $(seq 1 20); do
            sleep 5
            MY_HOSTNAME_NORM=$(echo "${MY_HOSTNAME:-}" | tr '_' '-')
            REPO_PREFIX="${MY_HOSTNAME_NORM%%-grocy-recipes}"
            for _c in "${REPO_PREFIX}-ha-storage" "${REPO_PREFIX}_ha_storage"; do
                if curl -sf --max-time 3 "http://${_c}:8099/api/products" >/dev/null 2>&1; then
                    STORAGE_URL="http://${_c}:8099"
                    break 2
                fi
            done
            if [ -z "${STORAGE_URL}" ] && [ -n "${SUPERVISOR_TOKEN}" ]; then
                _slug=$(curl -sf --max-time 5 \
                    -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
                    http://supervisor/addons \
                    | jq -r '[.data.addons[]? | select(.slug | test("ha.storage"; "i"))][0].slug // empty') || true
                if [ -n "${_slug}" ]; then
                    _ip=$(curl -sf --max-time 5 \
                        -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
                        "http://supervisor/addons/${_slug}/info" \
                        | jq -r '.data.ip_address // empty') || true
                    if [ -n "${_ip}" ]; then
                        STORAGE_URL="http://${_ip}:8099"
                        break
                    elif curl -sf --max-time 3 "http://${_slug}:8099/api/products" >/dev/null 2>&1; then
                        STORAGE_URL="http://${_slug}:8099"
                        break
                    fi
                fi
            fi
            bashio::log.info "Storage not found (attempt ${_attempt}/20) — retrying in 5 s..."
        done
        if [ -z "${STORAGE_URL}" ]; then
            bashio::log.warning "HA-Storage addon still not found after retries — falling back to localhost."
            STORAGE_URL="http://localhost:8099"
        else
            bashio::log.info "Storage URL (retry): ${STORAGE_URL}"
        fi
    fi
fi

# AI key/model are read directly from addon config
AI_PROVIDER="gemini"
if bashio::config.has_value 'ai_provider'; then
    AI_PROVIDER=$(bashio::config 'ai_provider')
fi

GEMINI_API_KEY=""
if bashio::config.has_value 'gemini_api_key'; then
    GEMINI_API_KEY=$(bashio::config 'gemini_api_key')
fi

GEMINI_MODEL="gemini-2.0-flash"
if bashio::config.has_value 'gemini_model'; then
    GEMINI_MODEL=$(bashio::config 'gemini_model')
fi

OLLAMA_URL=""
if bashio::config.has_value 'ollama_url'; then
    OLLAMA_URL=$(bashio::config 'ollama_url')
fi

OLLAMA_MODEL="llama3"
if bashio::config.has_value 'ollama_model'; then
    OLLAMA_MODEL=$(bashio::config 'ollama_model')
fi

CLAUDE_API_KEY=""
if bashio::config.has_value 'claude_api_key'; then
    CLAUDE_API_KEY=$(bashio::config 'claude_api_key')
fi

CLAUDE_MODEL="claude-3-5-haiku-20241022"
if bashio::config.has_value 'claude_model'; then
    CLAUDE_MODEL=$(bashio::config 'claude_model')
fi

DEBUG="false"
if bashio::config.true 'debug'; then
    DEBUG="true"
fi

export STORAGE_URL AI_PROVIDER GEMINI_API_KEY GEMINI_MODEL OLLAMA_URL OLLAMA_MODEL CLAUDE_API_KEY CLAUDE_MODEL DEBUG

exec python3 /app/backend.py
