# grocy_recipes/rootfs/etc/s6-overlay/s6-rc.d/recipe-backend/run
#!/command/with-contenv bashio

bashio::log.info "Starting recipe backend on port 8100..."

GROCY_BASE_URL=$(bashio::config 'grocy_base_url')
GROCY_API_KEY=$(bashio::config 'grocy_api_key')
GEMINI_API_KEY=$(bashio::config 'gemini_api_key')
GEMINI_MODEL=$(bashio::config 'gemini_model')

DEBUG="false"
if bashio::config.true 'debug'; then
    DEBUG="true"
fi

export GROCY_BASE_URL GROCY_API_KEY GEMINI_API_KEY GEMINI_MODEL DEBUG

exec python3 /app/backend.py
